# files/router_configs/common/etc/crontabs/root
#*/2 * * * * sh /etc/fabfi-scripts/check-uplink; sleep 30; sh /etc/fabfi-scripts/logstatus
*/2 * * * * sh /etc/fabfi-scripts/logstatus
*/1 * * * * sh /etc/fabfi-scripts/fabfimap update

