# branches/4.0.0_RC2/files/router_configs/wrt160nl-old/etc/crontabs/root
*/2 * * * * sh /etc/fabfi-scripts/check-uplink; sleep 30; sh /etc/fabfi-scripts/logstatus

